# files/desktop/etc/skel/.config/autostart/de_startup.desktop
[Desktop Entry]
Terminal=false
Type=Application
Name=User Desktop Environment Startup
Exec=$HOME/.local/bin/de_startup
