# files/general/usr/lib/udev/rules.d/98-mutuatech-performance-acplugged.rules
# disable power saving & enable performance settings when AC is plugged
SUBSYSTEM=="power_supply", ATTR{status}=="Charging", RUN+="/root/bin/Performance"
