# files/general/usr/lib/udev/rules.d/98-mutuatech-powersave-acunplugged.rules
# disable performance & enable power saving settings when AC is unplugged
SUBSYSTEM=="power_supply", ATTR{status}=="Discharging", RUN+="/root/bin/PowerSave"
