# files/server/etc/systemd/system/start-metatraders.timer
[Unit]
Description=Daily MetaTrader operation (start)

[Timer]
OnCalendar=Mon..Fri 09:45
Persistent=false

[Install]
WantedBy=timers.target
