# files/server/etc/systemd/system/start-metatraders.service
[Unit]
Description=timer: Starts (and forgets) local headless metatrader instances

[Service]
Type=simple
ExecStart=/root/bin/start-metatraders
KillMode=none
