# files/server/etc/systemd/system/stop-metatraders.timer
[Unit]
Description=Daily MetaTrader operation (stop)

[Timer]
OnCalendar=Mon..Fri 17:15
Persistent=false

[Install]
WantedBy=timers.target
